# .config/systemd/user/newsup.service
[Unit]
Description=Update newsboat

[Service]
Type=simple
ExecStart=/home/james/scripts/cron/newsup

[Install]
WantedBy=default.target
